# meta-rockM/recipes-apps/server/files/server.service
[Unit]
Description=server servie

[Service]
StandardOutput=null
StandardError=journal
Type=oneshot
SyslogIdentifier=server-service
ExecStart=/usr/bin/server

[Install]
WantedBy=multi-user.target
