# meta-rockM/recipes-apps/webapp/files/webapp.service
[Unit]
Description=server servie

[Service]
StandardOutput=null
StandardError=journal
Type=oneshot
SyslogIdentifier=server-service
ExecStart=/usr/bin/webrouter

[Install]
WantedBy=multi-user.target
